# linux/systemd/linapse-service.service
[Unit]
Description=CAD Mouse MK2 Linapse service (serial bridge + WebSocket)
After=ydotoold.service

[Service]
Environment=YDOTOOL_SOCKET=/run/user/%U/.ydotool_socket
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 %h/.local/bin/linapse-service
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
